# supervisord.conf
[unix_http_server]
file=/run/supervisor.sock
chmod=0770
chown=root:apache

[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/httpd -D FOREGROUND

[program:dnsmasq]
command=/usr/sbin/dnsmasq -k

[group:dnsmasq]
programs=dnsmasq

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface